# site-cookbooks/geminabox/templates/default/geminabox.socket.erb
[Unit]
Description=Gem in a box sockets

[Socket]
ListenStream=<%= @userhome %>/unicorn.sock
ListenStream=0.0.0.0:9191
Service=geminabox.service

[Install]
WantedBy=sockets.target
