# site-cookbooks/geminabox/templates/default/geminabox.socket.erb
[Unit]
Description=Gem in a box sockets

[Socket]
ListenStream=<%= @userhome %>/unicorn.sock
SocketGroup=www-data
Service=geminabox.service

[Install]
WantedBy=sockets.target
