# site-cookbooks/geminabox/templates/default/geminabox.service.erb
[Unit]
Description=Gem in a box

[Service]
ExecStart=<%= @clonedir %>/start.sh
WorkingDirectory=<%= @clonedir %>
KillSignal=SIGQUIT
Sockets=geminabox.socket
# weekly restart to force updates getting applied
RuntimeMaxSec=604800

[Install]
WantedBy=default.target
