# site-cookbooks/geminabox/templates/default/geminabox.service.erb
[Unit]
Description=Gem in a box

[Service]
ExecStart=<%= @clonedir %>/start.sh
WorkingDirectory=<%= @clonedir %>
KillSignal=SIGQUIT
Sockets=geminabox.socket

[Install]
WantedBy=default.target
