# site-cookbooks/zabbix/templates/default/zabbix-server.conf.erb
description "Zabbix Server"
author "Harald Sitter <sitter@kde.org>"

start on (filesystem and net-device-up IFACE!=lo)
stop on runlevel [!2345]

respawn
setuid zabbix
setgid zabbix

script
  exec /usr/sbin/zabbix_server
end script
