# site-cookbooks/publisher/templates/default/aptly_cleanup_upstart.conf.erb
author "Harald Sitter <sitter@kde.org>"
description "Cleanup Aptly database on system shutdown or reboot"

start on (runlevel [06]
          and (stopped <%= @uid %>_aptly
               and stopped <%= @uid %>_aptly_serve))

console log
setuid "<%= @user %>"
setgid "<%= @group %>"
chdir "<%= @dir %>"

script
  aptly db cleanup
  rm -rv aptly/upload/*
end script
