# site-cookbooks/publisher/templates/default/aptly.socket.erb
[Unit]
Description=Aptly API Socket

[Socket]
SocketMode=0660
ListenStream=<%= @listen_stream %>
SocketUser=<%= @user %>
SocketGroup=<%= @group %>

[Install]
WantedBy=sockets.target
