# site-cookbooks/publisher/templates/default/aptly_upstart.conf.erb
author "Rohan Garg <rohan@garg.io>"
description "Expose aptly API on publisher"

start on runlevel [2345]

console log
setuid "<%= @user %>"
setgid "<%= @group %>"
chdir "<%= @dir %>"

exec aptly api serve -listen "localhost:<%= @port %>"
