# site-cookbooks/publisher/templates/default/aptly.service.erb
[Unit]
Description=Aptly API

[Service]
ExecStart=/usr/bin/aptly api serve
Restart=always
WorkingDirectory=<%= @home %>

[Install]
WantedBy=default.target
