# site-cookbooks/gemstash/templates/default/gemstash.socket.erb
[Unit]
Description=Gemstash sockets

[Socket]
ListenStream=<%= @userhome %>/.gemstash/puma.pid
SocketGroup=www-data
Service=gemstash.service

[Install]
WantedBy=sockets.target
