# site-cookbooks/gemstash/templates/default/gemstash.service.erb
[Unit]
Description=Gemstash - A RubyGems.org cache and private gem server

[Service]
ExecStart=<%= @clonedir %>/start.sh
WorkingDirectory=<%= @clonedir %>
KillSignal=SIGQUIT
Sockets=geminabox.socket
# weekly restart to force updates getting applied
RuntimeMaxSec=604800

[Install]
WantedBy=default.target
